# peerberry/peerberry.service
[Unit]
Description=Execute 2fa

[Service]
Type=simple
ExecStart=/home/pi/2fa
Restart=on-failure
